# service/postamat-api.service
[Unit]
Description=Postamat API Service

[Service]
Type=simple
ExecStart= /root/postamat-api/run.sh

[Install]
WantedBy=multi-user.target